# deploy/cleanup-password-resets.service
[Unit]
Description=Cleanup expired/used password reset tokens for Tienda API
After=network.target

[Service]
Type=oneshot
User=rikashii
WorkingDirectory=/home/rikashii/tienda-api
ExecStart=/home/rikashii/tienda-api/.venv/bin/python /home/rikashii/tienda-api/scripts/cleanup_password_resets.py
EnvironmentFile=/home/rikashii/tienda-api/app/.env

[Install]
WantedBy=multi-user.target
